# templates/default/elasticsearch_exporter.upstart.erb
description "<%= @resource.service %>"

start on runlevel [2345]
stop on runlevel [016]

console log

exec <%= node.prometheus_client.install_dir%>/prometheus-consul-exporter -consul.server="<%= @resource.server %>" -web.listen-address=":<%= @resource.port %>"

respawn
respawn limit 3 180